# etc/init.d/amazon-ena
#!/sbin/openrc-run

depend() {
	return
}

start() {
	einfo "Enabling Amazon ENA support"
	modprobe -q ena
	eend 0
}
